# deploy/centos7/smart-inspect-backend.service
[Unit]
Description=Smart Inspect Agent FastAPI Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/smart-inspect-agent/backend
EnvironmentFile=/opt/smart-inspect-agent/.env
ExecStart=/opt/smart-inspect-agent/backend/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
